# Lab17Files/provision/playbooks/files/backupborg.service
[Unit]
Description=Backup creation script

[Service]
Type=simple
ExecStart=/usr/bin/borgbackup.sh

[Install]
WantedBy=default.target