# Lab12Files/provision/playbooks/files/catcher.service
[Unit]
Description=Monitor script

[Service]
Type=simple
ExecStart=/usr/bin/catcher.bash

[Install]
WantedBy=default.target